# deploy/geo-content.service
[Unit]
Description=GEO Content Platform API
After=network.target

[Service]
Type=simple
User=geo_content
Group=geo_content
WorkingDirectory=/opt/geo_content

# Environment variables for production
Environment=ENVIRONMENT=production
Environment=AWS_REGION=ap-southeast-1
Environment=AWS_SECRET_NAME=geo-content/api-keys
Environment=S3_BUCKET_UPLOADS=geo-content-uploads-053955129008
Environment=DB_PATH=/var/lib/geo_content/jobs.db
Environment=CORS_ORIGINS=http://geo-content-frontend-053955129008.s3-website-ap-southeast-1.amazonaws.com,http://18.138.26.110

# Use uv to run the application
ExecStart=/root/.local/bin/uv run uvicorn geo_content.main:app --host 127.0.0.1 --port 8000 --workers 2

# Restart configuration
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/geo_content/app.log
StandardError=append:/var/log/geo_content/error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
